# gdm/packages/gdm/postinst
#!/bin/bash

systemctl -f enable gdm
systemctl set-default graphical.target
usermod -G gdm,video gdm
